# clockspeed/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/clockspeed
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/clockadd
-rwxr-xr-x	root/root	usr/bin/clockfix
-rwxr-xr-x	root/root	usr/bin/clockspeed
-rwxr-xr-x	root/root	usr/bin/clockview
-rwxr-xr-x	root/root	usr/bin/sntpclock
-rwxr-xr-x	root/root	usr/bin/taiclock
-rwxr-xr-x	root/root	usr/bin/taiclockd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/clockadd.1.gz
-rw-r--r--	root/root	usr/share/man/man1/clockspeed.1.gz
-rw-r--r--	root/root	usr/share/man/man1/clockview.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sntpclock.1.gz
-rw-r--r--	root/root	usr/share/man/man1/taiclock.1.gz
-rw-r--r--	root/root	usr/share/man/man1/taiclockd.1.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/clockspeed/
drwxr-xr-x	root/root	var/lib/clockspeed/etc/
